# agent/container-build-files/build_files/etc/supervisord.d/15-lldpd.conf
[program:lldpd]
command=/usr/local/sbin/start-lldpd.sh
stdout_logfile=/var/log/supervisor/lldpd-daemon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
redirect_stderr=true