# fe/container-build-files/build_files/etc/supervisord.d/09-mariadb.conf
[program:mariadb]
command=bash -c "/usr/local/sbin/mariadb.sh"
stdout_logfile=/var/log/supervisor/siterm_mariadb.log
stdout_logfile_maxbytes=25MB
stdout_logfile_backups=2
redirect_stderr=true