# fe/container-build-files/build_files/etc/supervisord.d/09-mariadb.conf
[program:mariadb]
command=bash -c "/usr/local/sbin/mariadb.sh"
