# fe/container-build-files/build_files/etc/supervisord.d/12-httpd.conf
[program:httpd]
command=bash -c "sleep 10 && /usr/local/sbin/httpd.sh"
