# fe/container-build-files/build_files/etc/supervisord.d/12-httpd.conf
[program:httpd]
command=bash -c "sleep 10 && /usr/local/sbin/httpd.sh"
stdout_logfile=/var/log/supervisor/siterm_httpd.log
stdout_logfile_maxbytes=25MB
stdout_logfile_backups=2
redirect_stderr=true