# fe/container-build-files/build_files/etc/supervisord.d/14-provisioning-service.conf
[program:provisioning_service]
command=ProvisioningService-update --action start
priority=-1
stdout_logfile=/var/log/dtnrm-site-fe/ProvisioningService/daemon.log
stdout_logfile_maxbytes=0
redirect_stderr=true
umask=0002

[eventlistener:memmon_provisioning_service]
command=memmon -p provisioning_service=500MB
events=TICK_60
